# result/user/Library/LaunchAgents/com.bdsqqq.karabiner-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.bdsqqq.karabiner-build</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>/bin/wait4path /nix/store &amp;&amp; exec /nix/store/q0chpcz4qfwja47i4wv6hqh5dz79kjwy-karabiner-build-start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/karabiner-build.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/karabiner-build.log</string>
</dict>
</plist>